# toponavi/systemd/toponavi-pyfile@Location.service.d/exec.conf
[Service]
ExecStart=
ExecStart=/usr/bin/python3 -u PF_Location.py
